# BMS/battery.service
[Unit]
Description=Battery Monitoring Script

[Service]
WorkingDirectory=/path/to/Battery.py
ExecStart=/usr/bin/python3 /path/to/Battery.py
# Default to journald; setup.sh will overwrite with Desktop log files
StandardOutput=journal
StandardError=journal
Restart=always

[Install]
WantedBy=multi-user.target
